# scripts/systemd/phantom.service
[Unit]
Description=PHANTOM Raspberry Pi OS Application
Documentation=https://github.com/WilliamMajanja/Phantom
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=phantom
Group=phantom
WorkingDirectory=/opt/phantom
EnvironmentFile=-/etc/phantom/phantom.env
Environment=NODE_ENV=production
Environment=HOST=0.0.0.0
Environment=PORT=3000
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/phantom
ReadWritePaths=/etc/phantom

[Install]
WantedBy=multi-user.target
